# runit/log/run
#!/bin/bash
mkdir -p /sitelogs/drw/ratefood
exec svlogd /sitelogs/drw/ratefood
